# crontab
0 0 */2 * * cd /app && node index.js >> /proc/1/fd/1 2>&1